# smqueue/smqueue.conf
[program:smqueue]
command=sudo ./smqueue
directory=/OpenBTS/

